# vendor/unitedlinux/oracleasm.init.in
#! /bin/sh
# init fragment for asmlib
#
### BEGIN INIT INFO
# Provides: oracleasm
# Required-Start: hwscan
# Required-Stop:
# Default-Start: 2 3 5
# Default-Stop:
# Description: Load OCFS driver at system boot
### END INIT INFO

#
# Note that the start priority is after the network comes up.
# There is no real automatic mounting yet.
#

# Force LC_ALL=C
export LC_ALL=C

# Source init.d functions
. /etc/rc.status

prefix=@prefix@
exec_prefix=@exec_prefix@
libdir=@libdir@
. ${libdir}/oracleasm/oracleasm.init.functions

# Source networking configuration
[ -f /etc/sysconfig/oracleasm ] && . /etc/sysconfig/oracleasm


rc_reset

if [ -z "$ORACLE_ASMMANAGER" ]
then
    ORACLE_ASMMANAGER=/dev/oracleasm
fi
MODNAME=oracleasm

KVER="`uname -r`"
case "$KVER" in
2.4.19-64GB-SMP|2.4.19-4GB-SMP|2.4.19-4GB)
    MODPATH="/lib/modules/${KVER}/kernel/drivers/addon/oracleasm"
    ;;
2.4.21-*)
    BASEVER="2.4.21"
    REM="${KVER#${BASEVER}-}"
    EVER="${REM%%-*}"
    ETYPE="${REM#${EVER}-}"
    case "$ETYPE" in
    default|smp|psmp|numa|itanium2|itanium2-smp)
        ;;
    *)
        echo -n "Checking for United Linux Enterprise kernel:"
        false
        rc_status -v
        rc_exit
        ;;
    esac
    MODPATH="/lib/modules/${KVER}/kernel/drivers/addon/oracleasm"
    ;;
*)
    echo -n "Checking for United Linux Enterprise kernel:"
    false
    rc_status -v
    rc_exit
    ;;
esac




if_fail()
{
    RC="$1"
    if [ "$RC" = "0" ]
    then
        true
	rc_status -v
        return
    elif [ "$RC" = "2" ]
    then
        return
    fi
    false
    rc_status -v
    rc_exit
}


configure()
{
    configure_ask
    write_sysconfig
    if_fail $? "Unable to write the driver configuration"
}


scandisks()
{
    scan_disks
    if_fail "$?" "Unable to scan for ASM disks"
}


start()
{
    dev_create "${ORACLE_ASMMANAGER}"
    if_fail "$?" "Unable to create ${ORACLE_ASMMANAGER}"

    if [ "$ORACLEASM_ENABLED" != "true" -o ! -r "${MODPATH}/${MODNAME}.o" ]
    then
        return
    fi

    load_module "${MODNAME}"
    if_fail "$?" "Unable to load the ASMlib driver module"

    mount_device
    if_fail "$?" "Unable to mount the ASMlib driver filesystem"

    scandisks
}


stop()
{
    unmount_device
    if_fail "$?" "Unable to unmount the ASMlib driver filesystem"

    unload_module oracleasm
    if_fail "$?" "Unable to unload the ASM library driver"
}



case "$1" in
    start)
        start
	;;
	
    status)
        echo -n $"Checking if ASM is loaded: "
	if grep '^oracleasm ' /proc/modules >/dev/null 2>&1
        then
            rc_status -v
            echo -n $"Checking if $ORACLE_ASMMANAGER is mounted: "
            ORACLE_ASMMANAGERSEARCH="`echo $ORACLE_ASMMANAGER | sed -e 's/\//\\\\\//g'`"
            if grep "^oracleasmfs $ORACLE_ASMMANAGERSEARCH oracleasmfs" /proc/mounts >/dev/null 2>&1
            then
                rc_status -v
                rc_exit
            else
                rc_status -v
                rc_exit
            fi
        else
            false
            rc_status -v
            rc_exit
        fi
        ;;

    configure)
        configure
        if [ "$ORACLEASM_ENABLED" = "true" ]
        then
            start
        else
            stop
        fi
        ;;

    enable)
        ORACLEASM_ENABLED=true
        write_sysconfig
        if_fail $? "Unable to write the driver configuration"
        start
        ;;

    disable)
        ORACLEASM_ENABLED=false
        write_sysconfig
        if_fail $? "Unable to write the driver configuration"
        stop
        ;;

    createdisk)
        if [ -z "$2" -o -z "$3" ]
        then
            echo "Action \"createdisk\" requires two arguments" >&2
            exit 1
        fi
        make_disk "$2" "$3"
        if [ $? = 1 ]
        then
            if_fail 1 "Unable to create disk \"$2\""
        fi
        perm_disk "${ORACLE_ASMMANAGER}/disks/$2"
        if_fail $? "Unable to change ownership of disk \"$2\""
        ;;

    deletedisk)
        if [ -z "$2" ]
        then
            echo "Action \"deletedisk\" requires an argument" >&2
            exit 1
        fi
        remove_disk "$2"
        if_fail $? "Unable to delete disk \"$2\""
        ;;

    listdisks)
        ls -1 "${ORACLE_ASMMANAGER}/disks/"
        ;;

    querydisk)
        if [ -z "$2" ]
        then
            echo "Action \"querydisk\" requires an argument" >&2
            exit 1
        fi
        check_disk "$2"
        ;;

    scandisks)
        scandisks
        ;;

    stop)
        stop
        ;;

    restart)
        stop
        start
	;;

    *)
	echo "Usage: $0 {start|stop|restart|link|enable|disable|configure|createdisk|deletedisk|querydisk|listdisks|scandisks|status}"
        exit 1
esac

rc_exit
